# docker/supervisord.conf
# =====================================================
# LandIt Supervisor 配置
# 管理双服务：Nginx + Spring Boot
# @author Azir
# =====================================================

[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor
user=root

; 后端服务
[program:backend]
command=java -server -Xmx512m -Xms256m \
           -XX:+UseContainerSupport \
           -XX:MaxRAMPercentage=75.0 \
           -XX:InitialRAMPercentage=50.0 \
           -XX:+UseG1GC \
           -Djava.security.egd=file:/dev/./urandom \
           -jar /app/app.jar
directory=/app
user=app
autostart=true
autorestart=true
startretries=3
stderr_logfile=/var/log/supervisor/backend.err.log
stdout_logfile=/var/log/supervisor/backend.out.log
stopsignal=TERM
stopwaitsec=30
priority=1

; 前端服务
[program:nginx]
command=nginx -g 'daemon off;'
user=app
autostart=true
autorestart=true
startretries=3
stderr_logfile=/var/log/supervisor/nginx.err.log
stdout_logfile=/var/log/supervisor/nginx.out.log
stopsignal=QUIT
stopwaitsec=5
priority=2
